# fritzdect-bot.service
[Unit]
Description=FritzDECT Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=fritzdect
Group=fritzdect
WorkingDirectory=/opt/fritzdect-bot
ExecStart=/usr/bin/python3 /opt/fritzdect-bot/fritzdect_bot.py -c /opt/fritzdect-bot/config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fritzdect-bot

# Umgebungsvariablen
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
